# ops/systemd/kbo-discord-bot.service
[Unit]
Description=KBO Discord Bot
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/workspace/KBO_DiscordApp
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/workspace/KBO_DiscordApp/.env
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
